# assets/external/vacation.service
[Unit]
Description=Vacation Node.JS app

[Service]
User=ubuntu
WorkingDirectory=~/sharcnet-vacation
ExecStart=node .
Restart=always